# scripts/systemd/monitor_me_evidence_retention.timer
[Unit]
Description=Run MonitorMe evidence index scheduled retention hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=monitor_me_evidence_retention.service

[Install]
WantedBy=timers.target
